# ansible/roles/dns-master/files/named.service
[Unit]
Description=Named

[Service]
Type=forking
ExecStart=/bind/Docker/startup.script

[Install]
WantedBy=multi-user.target