# ansible/roles/dns-master/files/named.service
[Unit]
Description=Named

[Service]
Type=forking
ExecStart=ifconfig eth1:0 10.55.55.2 netmask 255.255.255.0 up

[Install]
WantedBy=multi-user.target